# systemd/user/mpd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/mpd --systemd %h/.config/mpd/mpd.conf
